# crontab
0 7 * * * /options=metadata/c/TugasKalbe/check_and_create_directory_main1.sh
